# Example Scripts/blacklist.timer
[Unit]
Description=Blacklist IP addresses from block list.

[Timer]
OnBootSec=15min
OnUnitActiveSec=23h

[Install]
WantedBy=timers.target
